# root/etc/cron.d/restart.cron
*/6 * * * * root /bin/sh /root/restart_wireguard.sh
@reboot  root  sleep 60 && /bin/sh /root/restart_wireguard.sh
